# x11-themes/app-misc/mythtv-status/files/mythtv-status.cron
#! /bin/sh

mythtv-update-motd --encoder-details --colour
